# truncate_all_tables.service
[Unit]
Description=Truncate all HCD and Presto tables for a fresh start
After=hcd.service


[Service]
Type=simple
Environment=PATH=/root/affiliate-junction-demo/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=VIRTUAL_ENV=/root/affiliate-junction-demo/.venv
WorkingDirectory=/root/affiliate-junction-demo
ExecStartPre=/bin/sleep 120
ExecStart=/root/affiliate-junction-demo/.venv/bin/python /root/affiliate-junction-demo/truncate_all_tables.py
Restart=on-failure
User=root
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
